# supervisord.conf
[supervisord] 

nodaemon=true 

 
 

[supervisorctl] 

serverurl=unix:///var/run/supervisor.sock 

 
 

[program:nginx] 

command=/usr/sbin/nginx -g 'daemon off;' 

autostart=true 

autorestart=true 

stdout_logfile=/dev/stdout 

stderr_logfile=/dev/stderr 

stdout_logfile_maxbytes=0 

stderr_logfile_maxbytes=0 

 
 

[program:python-service1] 

command=python /app/service1/app.py 

directory=/app/service1 

autostart=true 

autorestart=false 

stdout_logfile=/dev/stdout 

stderr_logfile=/dev/stderr 

stdout_logfile_maxbytes=0 

stderr_logfile_maxbytes=0 

# stderr_logfile=/var/log/supervisor/python-service1.err.log 

# stdout_logfile=/var/log/supervisor/python-service1.out.log 

 
 

[program:python-service2] 

command=python /app/service2/app.py 

directory=/app/service2 

autostart=true 

autorestart=true 

stdout_logfile=/dev/stdout 

stderr_logfile=/dev/stderr 

stdout_logfile_maxbytes=0 

stderr_logfile_maxbytes=0 

# stderr_logfile=/var/log/supervisor/python-service2.err.log 

# stdout_logfile=/var/log/supervisor/python-service2.out.log 

 
